# deploy/sibyl-builder/crontab
0 12 * * *     /usr/app/run.sh > /var/log/cron.log 2>&1

